# deploy/pin-web.service
[Unit]
Description=PIN Ticket Analysis web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/project/pin-web
EnvironmentFile=/home/ubuntu/project/pin-web/.env
Environment=PIN_WEB_HOST=0.0.0.0
Environment=PIN_WEB_PORT=8765
ExecStart=/home/ubuntu/project/pin-web/.venv/bin/python /home/ubuntu/project/pin-web/server.py
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20
StandardOutput=append:/home/ubuntu/project/pin-web/cache/pin-web.log
StandardError=append:/home/ubuntu/project/pin-web/cache/pin-web.log

[Install]
WantedBy=multi-user.target
